# overlay/etc/s6-overlay/s6-rc.d/wireguard/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-wireguard run ; run preflight"


#======================================================================================================================
# Start WireGuard interface
#======================================================================================================================

wg-quick up ${BF_WIREGUARD_CONF}
